# devops/php/assets/crontab
# min  	hour  	day	month 	week   	command
@reboot sleep 30 && echo "OK" >> /var/log/crontab.log
